# cron.d/syco-signer
*/5 * * * * root    /var/www/syco-signer/signer_trigger_delete.py
